# wificonnect/service/wificonnect.service
[Unit]
Description=WiFi Connect App Service
After=network-online.target

[Service]
User=root
WorkingDirectory=/root/wificonnect
ExecStart=/usr/bin/python3.8 /root/wificonnect/app2.py
Restart=always

[Install]
WantedBy=multi-user.target